# modulemanager/systemd/modulemanager.socket
[Socket]
ListenStream=12000

[Install]
WantedBy=sockets.target
